# .config/anacron/anacrontab
1         5 cron.daily   run-parts --report $XDG_CONFIG_HOME/anacron/cron.daily
7        10 cron.weekly  run-parts --report $XDG_CONFIG_HOME/anacron/cron.weekly
@monthly 15 cron.monthly run-parts --report $XDG_CONFIG_HOME/anacron/cron.monthly
